# cpp/sb.service
[Unit]
Description = Dry Run SB profiles

[Service]
Type=oneshot
ExecStart=sb-refresh --startup

[Install]
WantedBy=default.target
