# python/sb.service
[Unit]
Description = Dry Run SB profiles

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 1
ExecStart=sb-startup

[Install]
WantedBy=default.target
